# config/supervisor/1-xvfb.conf
[program:Xvfb]
command=Xvfb :0 -screen 0 %(ENV_RESOLUTION)sx24 -ac -nolisten unix
autorestart=true
priority=1
