# config/supervisor/4-websockify.conf
[program:Websockify]
command=websockify --web=/opt/novnc/ %(ENV_PORT)s localhost:5900
priority=4
autorestart=true
